# system/systemd-overrides/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin gort --noclear %I $TERM
